# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:backend]
# Listen on 0.0.0.0 so port 8000 is reachable directly from host
command=uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    UPLOAD_DIR="%(ENV_UPLOAD_DIR)s",
    SECRET_KEY="%(ENV_SECRET_KEY)s",
    FIRST_ADMIN_USERNAME="%(ENV_FIRST_ADMIN_USERNAME)s",
    FIRST_ADMIN_PASSWORD="%(ENV_FIRST_ADMIN_PASSWORD)s"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=2
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
